# django-uwsgi-nginx/supervisor-app.ini
[program:app-uwsgi]
command = /usr/bin/uwsgi --ini /home/docker/code/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx -c /home/docker/code/nginx-app.conf
